# vm/scripts/pve/cfg/rc.local
#!/bin/bash

#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.
 
# bash /root/bindip.sh

# 转发 k8s
ip route add 10.1.0.0/16 via 192.168.31.218 dev vmbr0 || true

exit 0